# files/etc/network/if-pre-up.d/ipset
#!/bin/sh
ipset restore < /etc/ipset.save
